# rc.local
#!/bin/sh -e
sudo date -s "$(wget -qSO- --max-redirect=0 google.com 2>&1 | grep Date: | cut -d' ' -f5-8)Z"
sudo mount -L HDCCK /raspi/MOUNTED_HD
echo $(date)>>/raspi/rc.local.dates.log
bash /usr/local/bin/raspi watcher >> /raspi/raspi.log &
exit 0
