# configs/raspbian/cron.d/mediacheky-update
@reboot carcheky /usr/local/bin/mediacheky update
@daily carcheky /usr/local/bin/mediacheky update
@reboot root bash /home/carcheky/mediacheky/raspiserver/scripts/bthelper.sh hci0
