# configs/raspbian/cron.d/mediacheky-update
@reboot carcheky /usr/local/bin/mediacheky update
@daily carcheky /usr/local/bin/mediacheky update
